# server-configs/supervisor-schema-sockets-demo.ini
[program:schemasocketsdemo]
command=/srv/http/schema-sockets-demo/demo/backend/server.py
directory=/srv/http/schema-sockets-demo/demo/backend
environment=PATH=/srv/http/schema-sockets-demo/env/bin
autorestart=true
user=nobody